# .d/sender.timer
# /etc/systemd/system/sender.timer

[Unit]
Description=Sender
Requires=sender.service
After=network-online.target

[Timer]
OnBootSec=10s
OnUnitActiveSec=10s
Persistent=true

[Install]
WantedBy=timers.target
